# systemd/weather-forecast.service
[Unit]
Description=Ruuvi Weather Forecast Orchestrator
Documentation=https://github.com/ruuvi/weather-forecast-system
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
EnvironmentFile=${RUUVI_INSTALL_DIR}/.env.weather
User=${RUUVI_SERVICE_USER}
Group=${RUUVI_SERVICE_USER}
WorkingDirectory=${RUUVI_INSTALL_DIR}
Environment=PATH=${RUUVI_INSTALL_DIR}/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=${RUUVI_INSTALL_DIR}
Environment=VIRTUAL_ENV=${RUUVI_INSTALL_DIR}/.venv
ExecStart=${RUUVI_INSTALL_DIR}/.venv/bin/python ${RUUVI_INSTALL_DIR}/scripts/weather_forecast_main.py --once
ExecReload=/bin/kill -HUP $MAINPID

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=${RUUVI_INSTALL_DIR}/logs ${RUUVI_INSTALL_DIR}/reports ${RUUVI_INSTALL_DIR}/data
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
RestrictNamespaces=true

# Restart policy
Restart=on-failure
RestartSec=30
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weather-forecast

# Timeout settings
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target